# OS/proc/os.manifest
# AIOS File Manifest
# Generated: 2026-04-04T03:49:36Z
# Format: sha256sum  path (relative to AIOS_ROOT)
#
# This manifest is used by system_check.sh to verify file integrity.
# Regenerate after updates: bash tools/generate-manifest.sh
#

12127e008d45f30874daea3ab6603f38e7f69b2476c70185ee093f9843b06afc  bin/aios-sys
15e4cbc0d2b929410c88ba6e3a698b66d057f2572afa03873fb81f265866f432  bin/aios
4b0760697886c9d7087dfed517b8613d634c78af48f8e3f34307392e4f446425  bin/aios-heartbeat
0fa07cf887801836951baa05e4e994882ec8211115dbfeb24313699609adb7e8  lib/aura-llama.sh
dec9efd80ab8554c7082046ac236d9d99073898697e4b8ff07f15ad7915ae091  lib/aura-security.sh
d5668e18ffdd7aa4509d0c4cd8848d00ca971842947dc845bf7461c8aaafb1a0  lib/aura-typo.sh
a0fb80cdd074a5e86bf8500cbe80a5d63cdf8f752792cceaef4396e3ca1ec7fe  lib/aura-proc.sh
02849da3d0c6ada2212470c3656edda2694db812dacc4175075e90845661258a  lib/aura-fs.sh
555e9864630797fb4e3cf555ece5a8aabeeae3c40e30ce92b4d369165f0ff97c  lib/aura-core.sh
7acc3f9625d06c467ab557915c08c81ff7c30d617d8b78418cd812d2439eda06  lib/aura-ai.sh
8566e8f2608053d5d86d19777df96140982e578354a0d461e529a69f09954c76  lib/aura-net.sh
6f0af3f6ddbc25a2d9fbe1659977aa587e26309bef42775a239d4901a2d832e4  tools/system_autofix.sh
1933ca2cfe9c3538db2b50f0b399e56d2048aba1614024fc06c508c32b1249fd  tools/log-viewer.sh
25ea0957b1990eca7c7e962d6c106bc8477aa9420affe6c061c4a2ff17ebfdc1  tools/offline-check.sh
a206403c0ea81ba7cc623724a7c6d952873a52aed4953f009b5c18750c61d615  tools/event-bus.sh
35401efd7883422d3d0fbf099737d459857c3f4c5211f672574d8fcb46ca21cd  tools/production-audit.sh
423c8b5004dc3a7bb3d2ba904279eb9e725799c951f1a21b348e91d3d2a2230e  tools/recovery-ctl.sh
cccd2a23918dc0f8ecebf0f2ee1729bb6aad5c1cd6330ce34d998e67f10c0c9d  tools/security-audit.sh
dc550cb23074286ed4179466e0335fd026574fe473786c9087b8001ae249b5b1  tools/perms-audit.sh
3fd56dda3263110f12a6dafe166cf60a60b2c2a0527dee2df8ef560936a83b86  tools/generate-manifest.sh
41b9f67756bf0379d80b49cd8b348607bb28fe12f4d58a0275029562fc344c87  tools/health_check.sh
4669cadbf1dee7e6ea04ebdcdf3d6b8a6c9deada6f2fd4cc9e6518245b9a9501  tools/service-ctl.sh
83037dea0fed4c1a5654780f995131ca16722c975f03ee20d94dfde1c09c65c4  tools/event-viewer.sh
952b3f7a0ea79b7116def5185349d0dd96e7ae6358cbb68ed93724f9a245711e  tools/ai-test.sh
7fa73014004522cafa0f4b887c4e21ac92e65caee1e4a34f91e54c9e6ae4d90b  tools/update-check.sh
bb5e34c8af0a278d44ad5500658b7d602152f9f4d6807c18c9caaa9b19a7f687  tools/detect-env.sh
d58e96b66cf24b093d8fc3c32c8e33eefccf24f39c827ec76fd2d29fdd5b53c5  tools/system_check.sh
50b7dda27a10bd292f3b83576e029d8a1240a7b49fe21fecbf2db53a110b790a  tools/module-ctl.sh
e8b629fae309f742f80bfc39fee3ecac6a064e19a789610d68c1467610b7cdb7  tools/config-audit.sh
5901e84a497e381c64d1a2706235184bd96abf6655955ab877dac0e617fb6a56  tools/apply-update.sh
738c3d1d0d244461e59eb4bda85b540880394cec3c963bee57f66617ca5771ac  config/network.conf
cf1c2bda705d2ff1c7a6e4b9d40bc34b69f2d193877fe94420fabe6530990694  config/api-endpoints.conf
6816b35291a1f3c05e7e2769660dd191b3d0f509b3191995cf8c6febc4818d3d  config/aios.conf
5e0758174f63fe84369d8945defcb557c05d59f750098d93768d8095fdd25132  config/device-profiles/termux.conf
837d5070c5404105cc76ae8a3af674abee472d5e761094ae4660db62896a777e  config/device-profiles/samsung-s21fe.conf
6e00674946fd5e9eea5b89db5bb3706762b956aef328700765a95468773263f7  config/device-profiles/generic-linux.conf
c6c695b7946536ce0f9881d9f2643650f8a918bb2275e011484c1417a3756873  config/services.conf
ed0ea8fa9aab148112e909b466eefb32ce018e3e172e1829a1c1c28363421987  config/system-manifest.conf
4bb6936e17589d18aafb348e0be69853b3944ce6d89e1b5462a69ca5974b505f  config/module-registry.conf
8e48b999a0ccd0b9ef8d9b1abf34f5bf36e0fafbe4945bbe3787521c839d4034  config/llama-settings.conf
7a67080aee23824e88cce9dfd22cfd7944e294d7e1c3db2e6704dabc650f0b98  etc/aios.conf
7100c5426f11ef78688f1dee530b8ce6826c6aa55cf1bc2b61c3b370f41f734d  OS/bin/shutdown
fde59cf3369c713e8975982b66fc8517733fec0abb48dcffbcf0d60a727f139b  OS/bin/os-state
c323ca70ffe80bb316b0b2949f81c0c0f6f5656fb44c678ae4c485d47e690aa5  OS/bin/os-netconf
dc91d8bbe32ba94e88f19638c41fae984d5f66462d4d225caf70e31faa1d821d  OS/bin/os-msg
9220b756cefd559b4652d09ad0fa37f845ce39b6ee8b266c3fae079f20b6f7d6  OS/bin/os-httpd
180bcea6e2bec03787fbc190e6ee5d7697d0de936f1db0059001142d1489ded0  OS/bin/os-event
d5be1b29e134340a0a96198066fb7e30683d4bbf76bdba17ba1befc5ceb8872d  OS/bin/os-check
4c92ccacd721a01f2b79a5e17a8d48571667f11439963a8e01e7d7f169f9085e  OS/bin/os-install
b5bd366890b2c14c402cfc100a3bf829c9a0efb31c9ffe65f71895fa6f67e78d  OS/bin/os-log
b20ace161a1bf337f3e0fc7c9337b1729a458f72d5bedb728423f3ea6e693236  OS/bin/os-service-status
8e19431011bdd784c90cc9883a47d984e4159d40ad0f12896c5b043b87b00a9d  OS/bin/login
de2ea1df6a9f5941a856375b4ea44080e47c367c2ec7bd1406ec4f9a18181100  OS/bin/os-checkpoint
220ffe65776c54b0f71c53877437f9b315abdf082f1619e0a6743fead385fcd7  OS/bin/os-service-health
1fdd9c3cbea29ededf786ea22c2de5de906021efa77567fab8ea55a142961e9a  OS/bin/os-health-wrapper
b6e11e6117c9e84a6a7b596937418b956d2c7d5c2c057ac9dd37ab17fe9835a7  OS/bin/os-svc
f372616cbcbbc59e7f1dc8d662b102745f49975781400cc7e3d814ba4b16b9fd  OS/bin/os-ai
32202550d0f2582635a59cb83b30f57eb3d1a9523a728437ecf76a89fd04f691  OS/bin/reboot
f827178b90d85afea760a826bfae6fa8a773dd27fa52a4d2227717356c272a30  OS/bin/sysinfo
35c0c49cbca757203d3b87d1135b9046e685337a671d1a8a1e7c17327931381c  OS/bin/os-resource
5097c01e7333d2dcb1a281d01518b7c7888f405877353159db5df3a3c8eb9368  OS/bin/os-syscall
61e50ba9f8dadc4fb9500631a70c6bec1342b044567921715a9abf64b9b12b91  OS/bin/busybox
ccecd8ebe8abcf10131d40746acfbc8c1183754f4178ae77dce8330f0c718958  OS/bin/os-perms
0a2308a815433766d91e4b6af71fe298f8e2f160d319a9b38752274b18578a5d  OS/bin/os-selftest
2d2dc00e8774c6aea576e94d965eb9b282c6a5177e2796250e61258554968277  OS/bin/os-login
8e9e1344948902ed18020487e37c7e0fe736d436640cdc46aca67d15c470b0af  OS/bin/os-recover
b397c25e9aef6930bd773115c122da1156b4b0deb8997ac9f696146bf1a0213b  OS/bin/os-ps
aacd669b4cc746a69eb179bcc5059a7b141576552a99e147d9ad737b8f802446  OS/bin/os-mirror
eb7da624b305c2a915ed4196dbeaf09d113a14b77e0fbfd5c1c96b5b8584722e  OS/bin/os-shell
3e8062e8fdcc8a734e70c53c805ff7c015d956a3c46aa8ca1c7939d012a11e04  OS/bin/os-service
a64761bebc0392ef8c362b7af1d5cb73a72eb7400846566b3631f5cd76331be6  OS/bin/os-real-shell
282de022914c3d7d4c837b15ced428c17ab84ad5980b28f298049d60bf7e9d03  OS/bin/os-sched
2696a326f0236c5d1be694cf7db725e735554005dd1f39010ad4f9974be12601  OS/bin/os-bridge
20173dc57f7c865261b3a76c3804400a58541b6d6ec9208e7f8b23c179e7bed5  OS/bin/os-kernelctl
183b395db0ef037de1a3cc1886b408b58ea505fce46a1a16916428d11a131363  OS/bin/os-emit
4fc30bad6af5921ff7d57b9a2d2fd1c89b0ddbebb1f395edf1b90210f8847b96  OS/bin/os-info
b0f07a44fba2aa60277bd756b26a2d9b77be86052c5caa3497e8cf430263a46c  OS/sbin/init.editbak
18ea5d4c3f192495ed3c621d456aef3449e1245435ed01d903c4ec04a3230be2  OS/sbin/init
aa52d6ff617d1693f8d447f2393618a4e50e39bfeed3c41be7931e4852889372  OS/etc/init.d/aura-bridge
33caa8b6b25e4aae90e500b51ed2753cd12ae0e4fd5064bbf5149355996e8b05  OS/etc/init.d/aura-tasks
584216bc47f5e745955deba9cd3e8235c1fd49ff60945f8fd8255b2515fb7615  OS/etc/init.d/devices
d61c68b6c1a0598edbe91b7ccb0f0cff1ce898503ff28ed7f21930e35eba3ef7  OS/etc/init.d/banner
797505213de6a54f23907b22a2f07ff3e6cec62657e51b1dee6252513259e324  OS/etc/init.d/os-kernel
f32adea50068d934f2ce3be3b5ee16b2e097e71b47016d9ad372f87532e981ba  OS/etc/init.d/aura-llm
a4bf9213478c289d53b4b3bb10d52d24f3867781da933e90b0b302efe080cfba  OS/etc/init.d/aura-agents
79125286ed663ac578c9d18b3f7591d0155022cf96d603d7faea9caccc17ed02  OS/etc/security.conf
d89683840f6a636042841fe0a4c763d78dafd005c5ab1003fbc982b9b303f971  ai/core/bots.py
49640551ef819820b075c224121677bebc56407a766d7eaf8f1824bd304031e4  ai/core/intent_engine.py
24649407367969486465831021173323e1f782eb2119c93c9b537b7761656df0  ai/core/commands.py
2a604a37a5acca64a514b67c9db1a8e6bd810be4135c6ef5d411ea9fa9725f4f  ai/core/router.py
453dcb9b714c19b5597f3d95ef3fee7244c40b2e6c5a3aa880c96804c0d6b462  ai/core/ai_backend.py
6e30e587e0efa37e7efe09cf13f30b3a0cdc68e2ea5ec796970b2d81bdd3d3de  ai/core/llama_client.py
42fab577ed178a4be8f13bf13467ef9eaf7e7f8a69ff91d6db3b513943dff2b2  ai/core/fuzzy.py
