# aioscpu/rootfs-overlay/systemd/system/getty@tty1.service.d/aioscpu-ai.conf
# AIOSCPU getty@tty1 drop-in override for AI mode
#
# This drop-in replaces the standard login prompt on tty1 with the
# AIOSCPU AI shell (auractl interactive) when aioscpu_mode=ai is set
# on the kernel command line.
#
# The aioscpu-mode-init.service writes the current mode to
# /run/aioscpu/mode so this condition can be evaluated.
#
# © 2026 Chris Betts | AIOSCPU Official | AI-generated, fully legal

[Unit]
# Only activate this override when the system booted in AI mode
ConditionPathExists=/run/aioscpu/mode
ConditionFileNotEmpty=/run/aioscpu/mode

[Service]
# Replace the standard agetty login with the AIOSCPU AI shell
ExecStart=
ExecStart=-/sbin/agetty --autologin aios --noclear %I $TERM
# After autologin the .bash_profile for aios will detect AI mode and
# exec aioscpu-ai-shell automatically (see rootfs-overlay profile).
